# deployment/com.feishu-chat.group-join.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.feishu-chat.group-join</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bytedance/Desktop/群聊消息统计/deployment/run_group_join.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/bytedance/Desktop/群聊消息统计</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- 每天北京时间 22:00 拉机器人进新群 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/bytedance/Desktop/群聊消息统计/logs/group-join.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/bytedance/Desktop/群聊消息统计/logs/group-join.err.log</string>
</dict>
</plist>
