# deployment/com.feishu-chat.daily-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.feishu-chat.daily-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bytedance/chorus-lark-monitor/deployment/run_daily_sync.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/bytedance/chorus-lark-monitor</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- 每天北京时间 10:00 同步 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/bytedance/chorus-lark-monitor/logs/daily-sync.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/bytedance/chorus-lark-monitor/logs/daily-sync.err.log</string>
</dict>
</plist>
