# deployment/com.feishu-chat.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.feishu-chat.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bytedance/chorus-lark-monitor/deployment/run_server.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/bytedance/chorus-lark-monitor</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <!-- 重启间隔（崩溃后），默认 10 秒，避免疯狂重启 -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/bytedance/chorus-lark-monitor/logs/server.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/bytedance/chorus-lark-monitor/logs/server.err.log</string>
</dict>
</plist>
